# etc/systemd/system/hostapd.service.d/wait-for-switch.conf
[Unit]
After=aic8800-switch.service
Requires=aic8800-switch.service
