# etc/systemd/system/aic8800-switch.service
[Unit]
Description=AIC8800D80 USB Modeswitch + Driver Loader
After=systemd-udevd.service
Requires=systemd-udevd.service

[Service]
Type=oneshot
RemainAfterExit=yes

# Stage 1: Modeswitch from 1111:1111 (CD-ROM) to a69c:8d80 (firmware loader)
# Device disconnects and re-enumerates after this
ExecStart=/usr/sbin/usb_modeswitch \
  -v 0x1111 -p 0x1111 \
  -M "555342438765432100000000000010fd0000000000000000000000000000f2"

# Wait for re-enumeration (stage 1->2)
ExecStartPost=/bin/sleep 10

# Stage 2: Load firmware downloader (triggers a69c:8d80 -> a69c:8d81 transition)
ExecStartPost=/sbin/modprobe aic_load_fw_usb

# Wait for firmware upload (stage 2->3)
ExecStartPost=/bin/sleep 5

# Stage 3: Load main WiFi driver (creates {{LAN_IFACE}})
ExecStartPost=/sbin/modprobe aic8800_fdrv_usb

# Wait for interface {{LAN_IFACE}} to appear (up to 30s, no brittle fixed sleep)
ExecStartPost=/bin/bash -c \
  'for i in $(seq 1 30); do \
     ip link show {{LAN_IFACE}} >/dev/null 2>&1 && break; \
     echo "Waiting for {{LAN_IFACE}}... ($i/30)"; \
     sleep 1; \
   done; \
   ip link show {{LAN_IFACE}} >/dev/null 2>&1 || \
     { echo "ERROR: {{LAN_IFACE}} did not appear after 30s"; exit 1; }'

# Assign static AP IP and bring interface up
# Placeholder {{AP_IP}}/{{AP_CIDR}} and {{LAN_IFACE}} are replaced by install.sh
ExecStartPost=-/sbin/ip addr add {{AP_IP}}/{{AP_CIDR}} dev {{LAN_IFACE}}
ExecStartPost=-/sbin/ip link set {{LAN_IFACE}} up
ExecStartPost=-/sbin/iw dev {{LAN_IFACE}} set power_save off

[Install]
WantedBy=multi-user.target
